# backend/deploy/supervisor/pme-conform-worker.conf
; Configuration Supervisor pour PME-CONFORM — workers Laravel.
;
; Ce fichier deploie 3 groupes de workers :
;
;   1. analyses       : analyse enrichie IA (longues, Ollama llama3.1:8b CPU).
;                       Timeout 4h pour absorber une analyse complete.
;
;   2. default        : jobs courts du systeme (notifications, exports XLSX,
;                       enrichissement post-analyse, etc.). Timeout standard.
;
;   3. referentiels   : indexation des referentiels uploades (extraction +
;                       chunks + embeddings nomic-embed-text + classification
;                       thematique LLM via ThemeClassifierService).
;
; Installation :
;   sudo cp pme-conform-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start pme-conform-analyses:* pme-conform-default:* pme-conform-referentiels:*
;
; Verification :
;   sudo supervisorctl status
;
; Redemarrage apres deploy (oblige pour charger le nouveau code) :
;   cd /home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel && php artisan queue:restart

; -----------------------------------------------------------------------
; GROUPE 1 — Worker analyses (enrichi IA, jobs longs)
; -----------------------------------------------------------------------
[program:pme-conform-analyses]
process_name=%(program_name)s_%(process_num)02d
command=php /home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel/artisan queue:work --queue=analyses --tries=1 --timeout=14400 --memory=2048 --sleep=3
directory=/home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=14400
user=asconsulting
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/pme-conform-analyses.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; -----------------------------------------------------------------------
; GROUPE 2 — Worker default (jobs courts, notifications, exports)
; -----------------------------------------------------------------------
[program:pme-conform-default]
process_name=%(program_name)s_%(process_num)02d
command=php /home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel/artisan queue:work --queue=default --tries=3 --timeout=300 --memory=512 --sleep=3
directory=/home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=300
user=asconsulting
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/pme-conform-default.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; -----------------------------------------------------------------------
; GROUPE 3 — Worker referentiels (indexation + classification thematique)
; -----------------------------------------------------------------------
[program:pme-conform-referentiels]
process_name=%(program_name)s_%(process_num)02d
command=php /home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel/artisan queue:work --queue=referentiels --tries=1 --timeout=3600 --memory=1024 --sleep=3
directory=/home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=3600
user=asconsulting
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/pme-conform-referentiels.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; -----------------------------------------------------------------------
; GROUPE 4 — Worker questionnaires (generation LLM des questionnaires DCP)
; -----------------------------------------------------------------------
; La generation d'un questionnaire IA appelle Ollama pour produire les
; questions par pole. Sur llama3.1:8b en CPU c'est 5-30 min selon la taille.
; Timeout 30 min pour absorber. Un seul worker suffit (rare qu'on genere
; plusieurs questionnaires en parallele).
[program:pme-conform-questionnaires]
process_name=%(program_name)s_%(process_num)02d
command=php /home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel/artisan queue:work --queue=questionnaires --tries=1 --timeout=1800 --memory=1024 --sleep=3
directory=/home/asconsulting/web/ns3485498.ip-193-70-32.eu/laravel
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=1800
user=asconsulting
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/pme-conform-questionnaires.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
